# services/systemd/actions-runner/actions-runner.service
[Unit]
Description=Actions Runner Service
After=network.target

[Service]
User=${USER}
WorkingDirectory=${PROJECT_DIR}/actions-runner
ExecStart=bash ${PROJECT_DIR}/actions-runner/run.sh
Restart=always
RestartSec=10
Environment=PATH=/usr/bin:/usr/local/bin:${PROJECT_DIR}/actions-runner
Environment=HOME=${HOME}
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

